# src/pkg/mac/postinstall
#!/bin/sh
set -e
/bin/launchctl load "/Library/LaunchDaemons/com.github.malliina.plist"